# deploy/ai.opencode.remote-control.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
                       "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.opencode.remote-control.telegram</string>

  <key>WorkingDirectory</key>
  <string>/Users/xtation/AgentWorks/Code_Opencode/opencode-remote-control</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>dist/launcher.js</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>NODE_ENV</key>
    <string>production</string>
  </dict>

  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>

  <key>StandardOutPath</key>
  <string>/tmp/opencode-remote-control-telegram.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/opencode-remote-control-telegram.err</string>

  <key>ThrottleInterval</key>
  <integer>10</integer>
</dict>
</plist>
